# rules/linux/process/at_command.yml
id: d2d642d7-b393-43fe-bae4-e81ed5915c4b
meta:
  author: phil
  description: Detects the use of at/atd.
  title: Scheduled Task/Job At
  references:
    - https://github.com/redcanaryco/atomic-red-team/blob/master/atomics/T1053.001/T1053.001.md
and: true
enabled: true
rules:
  exe_file:
    data:
      - /at
      - /atd
    type: endswith
source: process
system: linux
level: low
tags:
  - attack.persistence
  - attack.t1053.001