# rules/linux/process/lnx_schedule_task_job_cron.yml
id: 6b14bac8-3e3a-4324-8109-42f0546a347f
meta:
  author: phil
  description: Detects abuse of the cron utility to perform task scheduling for initial or recurring execution of malicious code. Detection will focus on crontab jobs uploaded from the tmp folder.
  title: Scheduled Cron Task/Job
  references:
    - https://github.com/redcanaryco/atomic-red-team/blob/master/atomics/T1053.003/T1053.003.md
and: true
enabled: true
rules:
  EventID:
    data:
      - '1002'
    type: string
  Image:
    data:
      - /crontab
    type: endswith
  Argv:
    data:
      - /tmp/
    type: contains
source: process
system: linux
level: medium
tags:
  - attack.execution
  - attack.persistence
  - attack.privilege_escalation
  - attack.t1053.003